# build/packages-template/bbb-web/imgprocess@.service
[Unit]
Description=Process an image file (Instance: %i)
After=network.target

[Service]
Type=oneshot
RemainAfterExit=no

TimeoutSec=60
MemoryHigh=512M
MemoryMax=640M
MemorySwapMax=0

ExecStart=/usr/local/bin/img-process.sh %i

[Install]
WantedBy=multi-user.target
