# record-and-playback/core/systemd/bbb_record_core.slice
[Unit]
Description = BigBlueButton Recording Processing
