# build/packages-template/bbb-config/bigbluebutton.target
[Unit]
Description=BigBlueButton System
Requires=
After=syslog.target network.target

[Install]
WantedBy=multi-user.target
