# build/packages-template/bbb-html5/disable-transparent-huge-pages.service
[Unit]
Description=Disable Transparent Huge Pages
ConditionPathIsDirectory=/sys/kernel/mm/transparent_hugepage

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/bin/echo "never" | tee /sys/kernel/mm/transparent_hugepage/enabled"
ExecStart=/bin/sh -c "/bin/echo "never" | tee /sys/kernel/mm/transparent_hugepage/defrag"
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
