# build/packages-template/bbb-webrtc-recorder/after-install.sh
#!/bin/bash -e

BIGBLUEBUTTON_USER=bigbluebutton

case "$1" in
  configure|upgrade|1|2)

    if id $BIGBLUEBUTTON_USER > /dev/null 2>&1 ; then
      chown $BIGBLUEBUTTON_USER:$BIGBLUEBUTTON_USER /var/lib/bbb-webrtc-recorder
      chmod 0700 /var/lib/bbb-webrtc-recorder
    fi

    # Set up the LiveKit API key/secret (generated by the bbb-livekit package)
    # as recorder env vars.
    if [ -f /etc/bigbluebutton/livekit.yaml ]; then
      # webhook.api_key is the active key; use it to pick the matching secret
      # rather than assuming a single/first entry under .keys.
      LIVEKIT_KEY=$(yq -r '.webhook.api_key' /etc/bigbluebutton/livekit.yaml)
      LIVEKIT_SECRET=$(yq -r ".keys.[\"$LIVEKIT_KEY\"]" /etc/bigbluebutton/livekit.yaml)

      if [ -n "$LIVEKIT_KEY" ] && [ "$LIVEKIT_KEY" != "null" ] \
        && [ -n "$LIVEKIT_SECRET" ] && [ "$LIVEKIT_SECRET" != "null" ]; then
        if [ ! -f /etc/default/bbb-webrtc-recorder ]; then
          touch /etc/default/bbb-webrtc-recorder
        fi
        chown $BIGBLUEBUTTON_USER:$BIGBLUEBUTTON_USER /etc/default/bbb-webrtc-recorder
        chmod 0640 /etc/default/bbb-webrtc-recorder

        # Re-sync the keys so changed keys in livekit.yaml propagate.
        sed -i '/^BBBRECORDER_LIVEKIT_APIKEY=/d; /^BBBRECORDER_LIVEKIT_APISECRET=/d' /etc/default/bbb-webrtc-recorder
        echo "BBBRECORDER_LIVEKIT_APIKEY=$LIVEKIT_KEY" >> /etc/default/bbb-webrtc-recorder
        echo "BBBRECORDER_LIVEKIT_APISECRET=$LIVEKIT_SECRET" >> /etc/default/bbb-webrtc-recorder
      elif [ -f /etc/default/bbb-webrtc-recorder ]; then
        # Lookup failed: drop stale creds so they don't drift from livekit.yaml.
        sed -i '/^BBBRECORDER_LIVEKIT_APIKEY=/d; /^BBBRECORDER_LIVEKIT_APISECRET=/d' /etc/default/bbb-webrtc-recorder
      fi
    fi

    systemctl enable bbb-webrtc-recorder
  ;;

  *)
    echo "## postinst called with unknown argument \`$1'" >&2
  ;;
esac

systemctl daemon-reload
